# supervisor.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi uwsgi.yml

[program:nginx]
command=nginx
